# mqtt-agent.service
[Unit]
Description=mqtt agent Service

[Service]
Type=simple
User=user
Group=user
WorkingDirectory=/PATH/TO/WORKINGDIRECTORY
ExecStart=/PATH/TO/MAIN/PROGRAM
ExecReload=/PATH/TO/MAIN/PROGRAM
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target